# .config/systemd/user/awww.service
[Unit]
Description=awww Wallpaper Daemon
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/awww-daemon
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
