# .config/systemd/user/hexstrike-server.service
[Unit]
Description=HexStrike AI - Security Tools API Server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/tools/hexstrike-ai
ExecStart=%h/tools/hexstrike-ai/hexstrike-env/bin/python3 %h/tools/hexstrike-ai/hexstrike_server.py --host 127.0.0.1 --port 8888
Restart=on-failure
RestartSec=5
Environment=HEXSTRIKE_HOST=127.0.0.1
Environment=HEXSTRIKE_PORT=8888

[Install]
WantedBy=default.target
