# etc/systemd/system/ollama.service.d/gpu.conf
[Service]
# Strix Halo (Radeon 8060S) iGPU offload via RADV/Vulkan.
# Requires AUR `ollama-vulkan-bin` (extra/ollama ships only CPU GGML libs).
Environment="OLLAMA_VULKAN=1"
Environment="OLLAMA_FLASH_ATTENTION=1"
Environment="OLLAMA_KV_CACHE_TYPE=q8_0"
Environment="OLLAMA_KEEP_ALIVE=30m"
